# rocrail-bin/files/rocrail.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2016 OpenWrt.org

START=99

start() {
	[ -e /opt/rocrail/bin ] && {
		cd /opt/rocrail/
		nohup ./bin/rocrail -nodevcheck -l /opt/rocrail/bin 2>&1 > /dev/null &
		echo "$!" > /tmp/rocraild.pid
	}
}

stop() {
	killall rocrail
}
